# etc/cron.hourly/backup_databases.sh
#!/bin/bash

cp /home/matt/studyjournal/db /home/matt/studyjournal/db.backup
cp /home/matt/language_study/db /home/matt/language_study/db.backup

chown matt:matt /home/matt/studyjournal/db.backup
chown matt:matt /home/matt/language_study/db.backup

