# etc/cron.hourly/ceden_backup_databases.sh
#!/bin/bash

cp /home/mg1/language_study/db /home/mg1/language_study/db.backup

chown mg1:mg1 /home/mg1/language_study/db.backup

